# deploy/crontab
# Ingesta Sinóptica (TZ del contenedor: America/Santiago)
# Observaciones METAR/DMC cada hora; pronósticos 2 veces al día; sismos cada
# 2 min (baja latencia de notificación) y estado de tsunami cada 10 min;
# incendios (FIRMS) y alertas (SENAPRED) cada hora; volcanes (RNVV) 2x/día;
# marea/oleaje (Open-Meteo Marine) 4x/día;
# infraestructura de emergencia (SENAPRED) y catastro de remociones en masa
# (SENAPRED) 1x/semana (ambas cuasi-estáticas); cortes de luz SEC y cortes de
# agua Esval cada 15 min y farmacias de turno MINSAL 2x/día (barato: solo leen
# data/incoming/* subido por el satélite en omen, sin red — ver
# satelite/README.md); precios
# de combustible CNE 2x/día (dormida sin CNE_API_KEY, accesible directo
# desde el VPS — no requiere el satélite); crecidas de ríos (GloFAS,
# Open-Meteo Flood API) 1x/día (el modelo global se actualiza diario).
# avisos meteo (derivados, propios) NO tienen línea propia: se recalculan
# dentro de las corridas de observaciones y de pronósticos (barato, sin red).
# watchdog de frescura (aviso a Slack) cada 10 min, dormido sin
# SLACK_WEBHOOK_URL — no es un paso de ingesta, se invoca directo como
# push/send.py.
4,14,24,34,44,54 * * * * python3 /app/ingesta/watchdog.py >> /data/ingesta.log 2>&1
# resumen diario de tráfico (Cloudflare Analytics) a Slack, dormido sin
# CF_ANALYTICS_TOKEN/CF_ZONE_ID/SLACK_WEBHOOK_URL — tampoco es paso de ingesta.
35 8 * * * python3 /app/ingesta/analytics_diario.py >> /data/ingesta.log 2>&1
7 * * * * python3 /app/ingesta/run.py >> /data/ingesta.log 2>&1
37 1,13 * * * python3 /app/ingesta/run.py --forecasts >> /data/ingesta.log 2>&1
*/2 * * * * python3 /app/ingesta/run.py --sismos >> /data/ingesta.log 2>&1
3-59/10 * * * * python3 /app/ingesta/run.py --tsunami >> /data/ingesta.log 2>&1
22 * * * * python3 /app/ingesta/run.py --incendios --alertas >> /data/ingesta.log 2>&1
52 8,20 * * * python3 /app/ingesta/run.py --volcanes >> /data/ingesta.log 2>&1
19 3,9,15,21 * * * python3 /app/ingesta/run.py --marea >> /data/ingesta.log 2>&1
15 5 * * * python3 /app/ingesta/run.py --crecidas >> /data/ingesta.log 2>&1
41 6 * * 1 python3 /app/ingesta/run.py --emergencia --remociones >> /data/ingesta.log 2>&1
*/15 * * * * python3 /app/ingesta/run.py --cortes >> /data/ingesta.log 2>&1
# esval desfasado 2 min de la grilla del satélite (sube a :00/:15/:30/:45):
# no leer el crudo mientras el scp lo está sobrescribiendo.
2-59/15 * * * * python3 /app/ingesta/run.py --esval >> /data/ingesta.log 2>&1
33 8,20 * * * python3 /app/ingesta/run.py --farmacias >> /data/ingesta.log 2>&1
44 7,19 * * * python3 /app/ingesta/run.py --combustible >> /data/ingesta.log 2>&1
# Mantener el log acotado (~100 KB)
12 4 * * * sh -c 'tail -c 100000 /data/ingesta.log > /data/ingesta.log.tmp && mv /data/ingesta.log.tmp /data/ingesta.log'
