# deploy/crontab
# Ingesta Sinóptica (TZ del contenedor: America/Santiago)
# Observaciones METAR/DMC cada hora; pronósticos 2 veces al día; sismos +
# estado de tsunami cada 10 min; incendios (FIRMS) y alertas (SENAPRED) cada
# hora; volcanes (RNVV) 2x/día; marea/oleaje (Open-Meteo Marine) 4x/día;
# infraestructura de emergencia (SENAPRED) y catastro de remociones en masa
# (SENAPRED) 1x/semana (ambas cuasi-estáticas).
# avisos meteo (derivados, propios) NO tienen línea propia: se recalculan
# dentro de las corridas de observaciones y de pronósticos (barato, sin red).
7 * * * * python3 /app/ingesta/run.py >> /data/ingesta.log 2>&1
37 1,13 * * * python3 /app/ingesta/run.py --forecasts >> /data/ingesta.log 2>&1
*/10 * * * * python3 /app/ingesta/run.py --sismos --tsunami >> /data/ingesta.log 2>&1
22 * * * * python3 /app/ingesta/run.py --incendios --alertas >> /data/ingesta.log 2>&1
52 8,20 * * * python3 /app/ingesta/run.py --volcanes >> /data/ingesta.log 2>&1
19 3,9,15,21 * * * python3 /app/ingesta/run.py --marea >> /data/ingesta.log 2>&1
41 6 * * 1 python3 /app/ingesta/run.py --emergencia --remociones >> /data/ingesta.log 2>&1
# Mantener el log acotado (~100 KB)
12 4 * * * sh -c 'tail -c 100000 /data/ingesta.log > /data/ingesta.log.tmp && mv /data/ingesta.log.tmp /data/ingesta.log'
