# deploy/crontab
# Ingesta Sinóptica (TZ del contenedor: America/Santiago)
# Observaciones METAR/DMC cada hora; pronósticos 2 veces al día; sismos cada
# 10 min; incendios (FIRMS) cada hora.
7 * * * * python3 /app/ingesta/run.py >> /data/ingesta.log 2>&1
37 1,13 * * * python3 /app/ingesta/run.py --forecasts >> /data/ingesta.log 2>&1
*/10 * * * * python3 /app/ingesta/run.py --sismos >> /data/ingesta.log 2>&1
22 * * * * python3 /app/ingesta/run.py --incendios >> /data/ingesta.log 2>&1
# Mantener el log acotado (~100 KB)
12 4 * * * sh -c 'tail -c 100000 /data/ingesta.log > /data/ingesta.log.tmp && mv /data/ingesta.log.tmp /data/ingesta.log'
